# website/supervisor/node.conf
[program:node]
command=node /app/.output/server/index.mjs
directory=/app
autostart=true
startsecs=10
autorestart=true
startretries=3
redirect_stderr=true
priority=10
