# website/supervisor/node.conf
[program:node]
command=node /app/.output/server/index.mjs
directory=/app
user=app
autostart=true
startsecs=10
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/app/logs/node.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=30
# Go server(:8000) 先起，Nuxt 反代 /api 才不会冷启动 502
priority=20
