# website/supervisor/server.conf
[program:server]
command=/app/server run -c /app/configs/config.toml
directory=/app
user=app
autostart=true
startsecs=10
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=30
stdout_logfile=/app/logs/server.log
# server 提供后端 API，需先于 Nuxt 启动
priority=10
